# crontab
MAILTO=""
52 * * * * /home/nobody/export1/k
53 * * * * /home/nobody/sta.sh
